# netwatch/data/netwatch.service
[Unit]
Description=NetWatch Network Traffic Analyzer
Documentation=https://github.com/netwatch
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python -m netwatch
WorkingDirectory=/opt/netwatch
User=root
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/netwatch /var/log/netwatch

# Allow network monitoring capabilities
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
